# files/etc/systemd/system/airplanes-webconfig.service
[Unit]
Description=airplanes.live webconfig
Documentation=https://github.com/airplanes-live/image
After=network.target lighttpd.service airplanes-first-run.service
Wants=network.target

[Service]
Type=simple
User=airplanes-webconfig
Group=airplanes-webconfig
# airplanes-feed: read access to /etc/airplanes/feeder-claim-secret (mode
# 0640 group=airplanes-feed) so the /api/identity/secret reveal handler can
# read the file directly without sudo. systemd loads supplementary groups
# via initgroups() at process start; this directive makes the dependency
# explicit instead of relying on /etc/group to be authoritative at unit
# load time.
SupplementaryGroups=airplanes-feed
ExecStart=/usr/local/bin/airplanes-webconfig --listen 127.0.0.1:8080
Restart=on-failure
RestartSec=5s

ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# /etc/airplanes is in ReadWritePaths so the sudo child running apply-config
# (as root) can write feed.env. Webconfig itself has no perms on
# /etc/airplanes (mode 0755 root:root) — only the helper, mediated by the
# sudoers entry, can mutate it. /run/airplanes is the tmpfiles-managed lock
# directory. /etc/NetworkManager/system-connections is where apl-wifi writes
# keyfiles; the sudo child inherits this unit's mount namespace so the path
# must be in ReadWritePaths or atomic-rename fails on the booted feeder
# (chroot tests don't see this).
ReadWritePaths=/var/lib/airplanes-webconfig /etc/airplanes /etc/airplanes/webconfig /run/airplanes /etc/NetworkManager/system-connections

# Webconfig invokes `sudo -n systemctl ...` for service restart, so
# NoNewPrivileges must stay false. The sudoers entry constrains the
# elevated argv to a hard-coded allowlist.
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
