# files/etc/systemd/system/airplanes-webconfig.service
[Unit]
Description=airplanes.live webconfig
Documentation=https://github.com/airplanes-live/image
After=network.target lighttpd.service airplanes-first-run.service
Wants=network.target

[Service]
Type=simple
User=airplanes-webconfig
Group=airplanes-webconfig
# airplanes-feed: read access to /etc/airplanes/feeder-claim-secret (mode
# 0640 group=airplanes-feed) so the /api/identity/secret reveal handler can
# read the file directly without sudo. video: /dev/vcio access for the
# hardware tile's `vcgencmd get_throttled` / `get_config psu_max_current`
# probes (mirrors the video-group grant feed makes for its own daemon user).
# systemd loads supplementary groups via initgroups() at process start; this
# directive makes the dependency explicit instead of relying on /etc/group to
# be authoritative at unit load time.
SupplementaryGroups=airplanes-feed video
# Session-token hashes are mirrored under /run so logins survive the
# service restart the runtime-overlay self-update performs (its rollback
# path stops the unit before restarting it, hence Preserve=yes rather
# than =restart). /run is tmpfs: a reboot clears all sessions. Mode 0700 —
# only the service account may read the mirror.
RuntimeDirectory=airplanes/webconfig
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
ExecStart=/opt/airplanes/current/bin/airplanes-webconfig --listen 127.0.0.1:8080 --sessions-path /run/airplanes/webconfig/sessions.json
Restart=on-failure
RestartSec=5s

ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
# /etc/airplanes is in ReadWritePaths so the sudo child running apply-config
# (as root) can write feed.env. Webconfig itself has no perms on
# /etc/airplanes (mode 0755 root:root) — only the helper, mediated by the
# sudoers entry, can mutate it. /run/airplanes is the tmpfiles-managed lock
# directory. /etc/NetworkManager/system-connections is where apl-wifi writes
# keyfiles; the sudo child inherits this unit's mount namespace so the path
# must be in ReadWritePaths or atomic-rename fails on the booted feeder
# (chroot tests don't see this).
#
# The apl-ssh helper (also a sudo child sharing this mount namespace under
# ProtectSystem=strict) writes the 99-airplanes-ssh-pi.conf snippet, the managed
# `pi` key file, and sets the pi password. chpasswd/passwd don't rewrite
# /etc/shadow in place — they create a temp file in /etc and rename over it, and
# take the /etc/.pwd.lock advisory lock — so /etc as a whole must be writable in
# this namespace, not just individual files. DAC still applies: the unprivileged
# service account owns none of these root-owned files, so only the sudo-root
# child can write them. Without /etc writable, chpasswd fails with EROFS on the
# booted feeder even though the sudoers grant authorises the call (chroot tests
# don't see this).
ReadWritePaths=/var/lib/airplanes/webconfig /etc/airplanes /etc/airplanes/webconfig /run/airplanes /etc/NetworkManager/system-connections /etc

# Webconfig invokes `sudo -n systemctl ...` for service restart, so
# NoNewPrivileges must stay false. The sudoers entry constrains the
# elevated argv to a hard-coded allowlist.
NoNewPrivileges=no

[Install]
WantedBy=multi-user.target
