# files/etc/systemd/system/airplanes-aggregator@.service
[Unit]
Description=airplanes.live third-party aggregator feed: %i
Documentation=https://airplanes.live
After=readsb.service network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/airplanes/current/lib/airplanes-webconfig/aggregator-run %i
Restart=always
RestartSec=10

# Vendor feeder code runs as a dedicated unprivileged account, not root. The
# account + state dir are provisioned by the image build (airplanes-live/image).
User=airplanes-aggregator
Group=airplanes-aggregator
# Only the writable state subtree is account-owned. The fetched vendor binary
# lives under the root-owned /var/lib/airplanes/aggregators/bin (read-only in
# this ProtectSystem=strict namespace), so the running account cannot replace
# the executable it runs.
StateDirectory=airplanes/aggregators/state
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
